# dist/fedora/nebula.service
[Unit]
Description=Nebula overlay networking tool

After=basic.target network.target network-online.target
Before=sshd.service
Wants=basic.target network-online.target

[Service]
WorkingDirectory=/etc/nebula
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/bin/nebula -config config.yml
Restart=always

[Install]
WantedBy=multi-user.target
